# rt-diff.cron
# Racktable Object diff
*/15 * * * * php /opt/rt-diff/bin/rt-diff.php
