# rt-diff.cron
# Racktable Object diff
*/15 * * * * root php /opt/rt-diff/bin/rt-diff.php
